# systemd/timer-unit/my-task.service
[Unit]
Description=My Task

[Service]
Type=oneshot
ExecStart=/path/to/your/script-or-command
# Replace with the actual command or script to execute

[Install]
WantedBy=multi-user.target
