# systemd/timer-unit/timer-unit.timer
[Unit]
Description=Run my task daily at 6 AM

[Timer]
# OnCalendar specifies when the timer should run
OnCalendar=daily
# Specify the time of day in UTC (24-hour format)
# Format: HH:MM:SS
OnCalendar=*-*-* 06:00:00
# Use persistent storage for missed runs
Persistent=true

[Install]
WantedBy=timers.target
